# lib/units/peach-update.timer
[Unit]
Description=Runs Peach Melpa hourly update

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target