# lib/units/peach-update.service
[Unit]
Description=Peach updater

[Service]
Type=simple
ExecStart=$HOME/peach-melpa/lib/units/peach-update.sh
